# infastructure_test/kafka/supervisord.conf
[supervisord]
nodaemon=true

[program:zookeper]
command=bin/zookeeper-server-start.sh config/zookeeper.properties

[program:kafka]
command=bin/kafka-server-start.sh config/server.properties